# docker/crontab
#!/bin/bash
# echo `date` > /var/www/html/storage/logs/schedule-run.log
* * * * * cd /var/www/html/ && php artisan schedule:run >> /dev/null 2>&1
